# systemd/antigravity-proxy-check.service
[Unit]
Description=Antigravity Proxy Login Check
After=network-online.target

[Service]
Type=oneshot
# 注意：安装时会自动替换 %SCRIPT_PATH% 为实际路径
# 或者手动修改为你的仓库路径，如：
# ExecStart=/path/to/your/repo/scripts/92_login_check.sh
ExecStart=%h/.local/bin/antigravity-login-check.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
